# ansible/files/90-monitor.rules
ACTION=="add", KERNEL=="wlan*", RUN+="/usr/bin/iwconfig %k mode monitor"
